# scripts/launchd/com.buddy.full.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.buddy.full</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/taehoonkwon/workspaces/buddy/core/scripts/buddy-core-full.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- NY 06:00 (EDT Local Time -> Sydney 20:00) -->
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>
</dict>
</plist>
